# etc/cron.d/gddu
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
*/15 * * * *   nobody    python3 /usr/local/bin/gddu.py
